# rootfs/var/lib/dpkg/info/apt.conffiles
/etc/apt/apt.conf.d/01autoremove
/etc/apt/apt.conf.d/20changelog
/etc/logrotate.d/apt
/etc/cron.daily/apt
